# .devcontainer/resources/etc/s6-overlay/s6-rc.d/dockerd-log/run
#!/bin/sh -e
exec logutil-service /var/log/dockerd