# .devcontainer/resources/etc/s6-overlay/s6-rc.d/dockerd-prepare-log/up
if { mkdir -p /var/log/dockerd }
if { chown nobody:nogroup /var/log/dockerd }
chmod 02755 /var/log/dockerd
